# offregister_etcd/data/etcd2.upstart.conf
description     "A highly-available key-value store for shared configuration and service discovery"
author          "Samuel Marks <samuelmarks@gmail.com>"

start on (net-device-up and local-filesystems and runlevel [2345])

stop on runlevel [016]

normal exit 0

respawn
respawn limit 10 5

env ETCD_ADVERTISE_CLIENT_URLS=http://%(public_ipv4)s:%(ADVERT_PORT)s
env ETCD_DISCOVERY=%(etcd_discovery)s
env ETCD_INITIAL_ADVERTISE_PEER_URLS=http://%(public_ipv4)s:%(PEER_PORT)s
env ETCD_LISTEN_CLIENT_URLS=http://0.0.0.0:%(ADVERT_PORT)s,http://0.0.0.0:%(ADDITIONAL_LISTEN_PORT)s
env ETCD_LISTEN_PEER_URLS=http://0.0.0.0:%(PEER_PORT)s
env ETCD_NAME=%(node_name)s
env ETCD_DATADIR=%(data_dir)s

exec etcd
